# crates/gateway/gateway.service
[Unit]
Description=agent-gateway
Documentation=https://github.com/nitecon/agent-gateway
After=network.target

[Service]
Type=simple
User=agentic
Group=agentic

EnvironmentFile=/etc/agent-gateway/gateway.env

# Deploys are owned by Eventic (see .eventic.yaml + deploy/eventic-deploy.sh).
# The unit just runs whatever binary the symlink currently resolves to — no
# self-update on start.
ExecStart=/usr/local/bin/gateway

Restart=on-failure
RestartSec=5
TimeoutStopSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-gateway

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/agentic

[Install]
WantedBy=multi-user.target
