# crates/gateway/gateway.service
[Unit]
Description=agent-gateway
Documentation=https://github.com/nitecon/agent-gateway
After=network.target

[Service]
Type=simple
User=agentic
Group=agentic

# Secrets and config live in /etc/agent-gateway/gateway.env (not world-readable).
# Copy crates/gateway/.env.example to /etc/agent-gateway/gateway.env and fill in values.
EnvironmentFile=/etc/agent-gateway/gateway.env

ExecStartPre=/usr/local/bin/gateway update
ExecStart=/usr/local/bin/gateway

Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Logging goes to the system journal: journalctl -u gateway -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-gateway

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/agent-gateway /opt/agentic/bin

[Install]
WantedBy=multi-user.target
